# etc/supervisord.conf
[program:wadowedo]
command=/home/ubuntu/.nvm/current/bin/node --harmony main.js
directory=/home/ubuntu/wadowedo/wadowedo/
user=ubuntu
environment=NODE_ENV=staging