# etc/supervisord.conf
[program:wadowedo]
command=/usr/local/bin/node --harmony main.js
directory=/home/ubuntu/wadowedo
user=root
environment=NODE_ENV=staging